# 98-ipchrg.rules
ENV{DEVTYPE}=="usb_device", ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="05ac", ATTR{idProduct}=="12a8", RUN+="/usr/bin/ipchrg"
